# supervisor/redis.ini
[program:redis-master]
command=/bin/su redis /bin/bash -c "redis-server /opt/redis/conf/master.conf"
process_name=%(program_name)s
pidfile=/var/run/redis-master.pid
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:redis-replica]
command=/bin/su redis /bin/bash -c "redis-server /opt/redis/conf/replica.conf"
process_name=%(program_name)s
pidfile=/var/run/redis-replica.pid
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
